# ansible/roles/pico-web/templates/ctf-stats.timer.j2
[Unit]
Description=ctf stats cache updater (timer)

[Timer]
OnUnitInactiveSec=60s
OnBootSec=30s

[Install]
WantedBy=timers.target
